# Quadlet/vpn-wireguard-container.container
[Container]
ContainerName=Container-vpn-wireguard
Image=localhost/vpn
Pod=VPN.pod
DropCapability=all
AddCapability=CAP_NET_ADMIN
AutoUpdate=local
